# keepalive/crontab
# m  h  dom mon dow   command
*/2  *   *   *   *    sh /home/pi/meshtastic-live.sh
0    0   *   *   *    rm /home/pi/cronlog
